# filesystem/etc/skel/.config/autostart/guake.desktop
[Desktop Entry]
Type=Application
Exec=guake
Hidden=false
X-MATE-Autostart-enabled=true
Name[en_US]=Guake Terminal
Name=Guake Terminal
Comment[en_US]=
Comment=
